# helper_scripts/systemd/openclaw-tls-renew.service
# MODULE_NOTE (CN)：
#   OPS-1 Track A：Tailscale cert renewal one-shot service。
#   為什麼 Type=oneshot：跑完即退；由 sibling timer 每日 03:00 UTC 觸發。
#   為什麼 OnFailure：spec §2.2 — renewal 失敗時觸 notify script（Telegram +
#   Grafana alert）。
#
# 環境：
#   OPENCLAW_TLS_CERT_DIR / OPENCLAW_TLS_CERT_HOST 由 Environment= 傳入
#   或讓 helper script 讀 secrets env file。

[Unit]
Description=OpenClaw Tailscale cert renewal (one-shot, fired by timer)
After=tailscaled.service
Wants=tailscaled.service
OnFailure=openclaw-tls-renew-notify.service

[Service]
Type=oneshot
# 必要 env：可由 systemctl edit override 注入
Environment="OPENCLAW_BASE_DIR=/opt/openclaw"
# 為什麼用 helper script：跨平台抽象 + 14d threshold + reload Caddy
ExecStart=/bin/bash -c '\
    source ${OPENCLAW_BASE_DIR}/helper_scripts/lib/tls_cert.sh; \
    host="$(resolve_openclaw_tls_cert_host)"; \
    dir="$(resolve_openclaw_tls_cert_dir)"; \
    if tls_cert_should_renew "$dir/$host.crt" 14; then \
        tailscale cert "$host" && \
        chown root:caddy "$dir/$host.crt" "$dir/$host.key" && \
        chmod 640 "$dir/$host.crt" "$dir/$host.key" && \
        systemctl reload openclaw-caddy.service; \
    else \
        echo "cert still has >14d remaining; skip renewal"; \
    fi'

[Install]
# 通常不直接 enable；由 sibling timer 觸發
WantedBy=multi-user.target
