# helper_scripts/systemd/openclaw-listing-collector.service
# MODULE_NOTE
# 模塊用途：Linux systemd unit — production listing capture-only collector 常駐 daemon
#   （COLLECTOR-LISTING-CAPTURE-PROD；OQ-2 systemd 常駐）。前向累積新上市 perp 上市
#   瞬間逐筆 publicTrade + phase transition + capture_lag，雙寫 market.klines（additive）
#   + research.listing_capture_events。
# 環境占位符（install_listing_collector_service.sh 安裝時 sed 替換）：
#   __COLLECTOR_USER__ / __COLLECTOR_GROUP__ — 跑 collector 的 OS 帳號（通常 ncyu）
#   __OPENCLAW_BASE_DIR__ — srv repo root（如 /home/ncyu/BybitOpenClaw/srv）
#   __OPENCLAW_DATA_DIR__ — runtime data dir（JSONL fallback + lock 落此）
#   __OPENCLAW_SECRETS_ROOT__ — secrets 根（PG 連線參數來源）
#   __PYTHON_BIN__ — python3 解譯器路徑（如 /usr/bin/python3）
# 硬邊界（capture-only 旁路）：
#   1. capture-only：零 order / 零 strategy intent / 零 IPC trading / 零 live /
#      零 execution_authority / 零 改 engine。寫 PG 限 research 表 + market.klines additive。
#   2. 不需 OPENCLAW_ALLOW_MAINNET / authorization.json / IPC socket（純資料採集 daemon）。
#   3. Restart=always（與 engine 的 on-failure 不同）：collector 是純資料採集，無 live-auth
#      風險；listing 不可重捕，daemon 任何退出（含 SIGTERM exit-0）都應重啟續捕，最大化
#      捕獲連續性。operator 主動停用走 systemctl stop + disable（不靠 exit code 區分）。
#   4. StartLimitBurst=10 / StartLimitIntervalSec=300：10 連 fail / 5min 內進 failed
#      防 thrash（WS/PG 全掛時不無限重啟燒 CPU）。
#   5. PreStart fail-closed：collector daemon 入口 + python 存在 + secrets env 存在。
#
# 與 engine unit 的差異（為什麼不共用）：
#   - Restart=always vs engine on-failure（capture 連續性 > 主動 stop 語義）。
#   - 無 IPC socket / 無 mainnet env / 無 binary（純 Python daemon）。
#   - 無 5-gate 關聯（collector 不碰 authorization.json）。
#
# 安裝步驟：
#   sudo bash helper_scripts/systemd/install_listing_collector_service.sh
#
# 啟停命令：
#   sudo systemctl start openclaw-listing-collector
#   sudo systemctl stop openclaw-listing-collector
#   sudo systemctl status openclaw-listing-collector
#   sudo journalctl -u openclaw-listing-collector -f

[Unit]
Description=Arcane Equilibrium / OpenClaw Listing Capture-Only Collector
After=network-online.target postgresql.service
Wants=network-online.target
# PG 不在同機亦不阻 collector 啟動（collector 內含 PG retry + JSONL fallback）
# 不設 Requires= — 故意無硬依賴

# systemd v245+ 將 StartLimit* 視為 [Unit] directive
StartLimitBurst=10
StartLimitIntervalSec=300

[Service]
Type=simple
User=__COLLECTOR_USER__
Group=__COLLECTOR_GROUP__
WorkingDirectory=__OPENCLAW_BASE_DIR__

# EnvironmentFile — 載入 basic_system_services.env（含 POSTGRES_* PG 連線參數）
EnvironmentFile=__OPENCLAW_SECRETS_ROOT__/environment_files/basic_system_services.env

# 補上 base/data/secrets dir（pg_sink read_db_config + JSONL fallback 用）
Environment="OPENCLAW_BASE_DIR=__OPENCLAW_BASE_DIR__"
Environment="OPENCLAW_DATA_DIR=__OPENCLAW_DATA_DIR__"
Environment="OPENCLAW_SECRETS_ROOT=__OPENCLAW_SECRETS_ROOT__"

# PreStart fail-closed 三檢：collector daemon 入口 / python / secrets env 存在
ExecStartPre=/usr/bin/test -f __OPENCLAW_BASE_DIR__/helper_scripts/collectors/listing_capture/daemon.py
ExecStartPre=/usr/bin/test -x __PYTHON_BIN__
ExecStartPre=/usr/bin/test -s __OPENCLAW_SECRETS_ROOT__/environment_files/basic_system_services.env

# collector daemon 啟動（常駐，無 --smoke-seconds）
ExecStart=__PYTHON_BIN__ __OPENCLAW_BASE_DIR__/helper_scripts/collectors/listing_capture/daemon.py

# Graceful stop：daemon 有 SIGTERM → stop event → 關 WS + flush PG + exit
KillSignal=SIGTERM
TimeoutStopSec=20

# 重啟策略 — capture 連續性優先
# Restart=always：listing 不可重捕，任何退出都重啟續捕（含 SIGTERM exit-0）
# RestartSec=10：避免 race + 給 WS/PG 短暫恢復時間
Restart=always
RestartSec=10

# 日誌 — collector 用 logging.basicConfig 寫 stdout/stderr；systemd capture 到 journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-listing-collector

# 資源限制 — 避 fork-bomb / fd leak
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
