# helper_scripts/deploy/com.openclaw.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  MODULE_NOTE (中文):
    macOS launchd 配置模板 — OpenClaw Gateway 服务（Node.js）。
    对应 systemd 的 openclaw-gateway.service，功能完全等价。
    AI 供应商密钥必须通过环境变量传入，禁止硬编码。

  MODULE_NOTE (English):
    macOS launchd plist template — OpenClaw Gateway service (Node.js).
    Equivalent to the systemd openclaw-gateway.service unit.
    AI provider keys MUST be passed via environment variables; never hardcode.

  安裝步驟 / Installation:
    1. cp com.openclaw.gateway.plist ~/Library/LaunchAgents/
    2. sed -i '' "s|__HOME__|$HOME|g" \
           ~/Library/LaunchAgents/com.openclaw.gateway.plist
       # Gateway 不參考 repo 根路徑（Node 模組裝在 npm-global），
       # 所以無需 __BASE__ sed（保留此行僅為形狀一致性）。
    3. 透過 launchctl setenv 或密钥管理器注入 AI 密钥；不要编辑 plist 写入密钥。
       Inject AI keys with launchctl setenv or a secrets manager; do not edit keys into this plist.
    4. mkdir -p ~/Library/Logs/openclaw
    5. launchctl load ~/Library/LaunchAgents/com.openclaw.gateway.plist

  安全提醒 / Security Warning:
    不要在 plist 中保存 OPENAI_API_KEY / ANTHROPIC_API_KEY。
    Do NOT store OPENAI_API_KEY / ANTHROPIC_API_KEY in this plist.
-->
<plist version="1.0">
<dict>

  <!-- 服务标识 / Service identifier -->
  <key>Label</key>
  <string>com.openclaw.gateway</string>

  <!-- 启动命令 / Start command
       对应 systemd ExecStart: /usr/bin/node ...openclaw/dist/entry.js gateway (port 18789)
       macOS 上 node 路径可能不同，需根据实际安装调整
       Node path may differ on macOS; adjust to actual installation -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>__HOME__/.npm-global/lib/node_modules/openclaw/dist/entry.js</string>
    <string>gateway</string>
    <string>--port</string>
    <string>18789</string>
  </array>

  <!-- 自动重启 / Auto-restart on crash -->
  <key>KeepAlive</key>
  <true/>

  <!-- 重启间隔 / Restart throttle -->
  <key>ThrottleInterval</key>
  <integer>5</integer>

  <!-- 停止超时 / Stop timeout -->
  <key>ExitTimeOut</key>
  <integer>30</integer>

  <!-- 环境变量 / Environment variables
       安全提醒：此模板只包含非敏感运行元数据；供应商密钥放 OS Keychain、
       secret manager 或 0600 文件，由服务 wrapper 读取，不写入 plist / launchd env。
       Security: this template contains only non-sensitive runtime metadata;
       provider keys belong in OS Keychain, a secret manager, or 0600 files. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>

    <key>TMPDIR</key>
    <string>/tmp</string>

    <key>PATH</key>
    <string>__HOME__/.local/bin:__HOME__/.npm-global/bin:/usr/local/bin:/usr/bin:/bin</string>

    <!-- Gateway 服务元数据 / Service metadata -->
    <key>OPENCLAW_GATEWAY_PORT</key>
    <string>18789</string>

    <key>OPENCLAW_SERVICE_MARKER</key>
    <string>openclaw</string>

    <key>OPENCLAW_SERVICE_KIND</key>
    <string>gateway</string>

    <key>OPENCLAW_SERVICE_VERSION</key>
    <string>2026.4.1</string>
  </dict>

  <!-- 日志文件 / Log files -->
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/openclaw/gateway-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/openclaw/gateway-stderr.log</string>

  <!-- 登录时自动运行 / Run at load -->
  <key>RunAtLoad</key>
  <true/>

</dict>
</plist>
