# helper_scripts/deploy/com.openclaw.trading-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  MODULE_NOTE (中文):
    macOS launchd 配置模板 — OpenClaw Trading API 服务。
    对应 systemd 的 openclaw-trading-api.service，功能完全等价。
    使用 __BASE__（repo 根）+ __HOME__（家目錄）雙占位符，安裝時 sed 替換。
    repo 若部署於非 $HOME 路徑（如 /Users/ncyu/Documents/Projects/TradeBot），
    __BASE__ 可設為該絕對路徑。放置於 ~/Library/LaunchAgents/ 實現開機自啟。

  MODULE_NOTE (English):
    macOS launchd plist template — OpenClaw Trading API service.
    Equivalent to the systemd openclaw-trading-api.service unit.
    Uses __BASE__ (repo root) + __HOME__ (user home) placeholders.
    Place in ~/Library/LaunchAgents/ for auto-start on login.

  安裝步驟 / Installation:
    1. cp com.openclaw.trading-api.plist ~/Library/LaunchAgents/
    2. sed -i '' "s|__BASE__|$OPENCLAW_BASE_DIR|g" \
           ~/Library/LaunchAgents/com.openclaw.trading-api.plist
       sed -i '' "s|__HOME__|$HOME|g" \
           ~/Library/LaunchAgents/com.openclaw.trading-api.plist
    3. mkdir -p ~/Library/Logs/openclaw
    4. launchctl load ~/Library/LaunchAgents/com.openclaw.trading-api.plist

  对应 systemd 配置映射 / systemd Mapping:
    WorkingDirectory  → WorkingDirectory key
    ExecStart         → ProgramArguments array
    Restart=always    → KeepAlive=true
    RestartSec=5      → ThrottleInterval=5
    TimeoutStopSec=30 → ExitTimeOut=30
    Environment=      → EnvironmentVariables dict
    WantedBy=default  → placed in ~/Library/LaunchAgents/ (auto-load on login)
-->
<plist version="1.0">
<dict>

  <!-- 服务标识 / Service identifier -->
  <key>Label</key>
  <string>com.openclaw.trading-api</string>

  <!-- 工作目录 / Working directory
       对应 systemd: WorkingDirectory=
       Equivalent to systemd WorkingDirectory= -->
  <key>WorkingDirectory</key>
  <string>__BASE__/program_code/exchange_connectors/bybit_connector/control_api_v1</string>

  <!-- 启动命令 / Start command
       对应 systemd: ExecStart=
       uvicorn 通过 venv 内的完整路径调用，确保使用正确的 Python 环境
       Uses full venv path to ensure correct Python environment -->
  <key>ProgramArguments</key>
  <array>
    <string>__BASE__/program_code/exchange_connectors/bybit_connector/control_api_v1/.venv/bin/uvicorn</string>
    <string>app.main:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
  </array>

  <!-- 自动重启 / Auto-restart on crash
       对应 systemd: Restart=always
       KeepAlive=true 确保进程异常退出后自动重启 -->
  <key>KeepAlive</key>
  <true/>

  <!-- 重启间隔（秒）/ Restart throttle interval
       对应 systemd: RestartSec=5
       防止崩溃循环时过于频繁重启 -->
  <key>ThrottleInterval</key>
  <integer>5</integer>

  <!-- 停止超时（秒）/ Stop timeout
       对应 systemd: TimeoutStopSec=30
       超时后 launchd 发送 SIGKILL -->
  <key>ExitTimeOut</key>
  <integer>30</integer>

  <!-- 环境变量 / Environment variables
       对应 systemd: Environment= 行
       仅包含必需的运行时变量；密钥通过文件注入，不在此硬编码
       Only runtime-required vars; secrets injected via files, never hardcoded here -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME__</string>

    <key>TMPDIR</key>
    <string>/tmp</string>

    <key>PATH</key>
    <string>__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin</string>

    <!-- 时区 / Timezone -->
    <key>TZ</key>
    <string>Europe/Madrid</string>

    <!-- Bybit API 目标（testnet/mainnet）/ Bybit API target
         testnet: https://api-testnet.bybit.com
         mainnet: https://api.bybit.com -->
    <key>BYBIT_API_HOST</key>
    <string>https://api-testnet.bybit.com</string>
  </dict>

  <!-- 标准输出日志 / Stdout log
       macOS 无 journalctl，日志写入文件
       macOS has no journalctl; logs go to files -->
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/openclaw/trading-api-stdout.log</string>

  <!-- 标准错误日志 / Stderr log -->
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/openclaw/trading-api-stderr.log</string>

  <!-- 登录时自动运行 / Run at load (on login)
       对应 systemd WantedBy=default.target + enable -->
  <key>RunAtLoad</key>
  <true/>

</dict>
</plist>
